# roles/sys-net/templates/etc/systemd/system/networking.service.d/override.conf.j2
# {{ ansible_managed }}
# ansibleguy.linux_networking

[Service]
ExecStartPost=/usr/bin/systemctl restart {{ nw.startup_service.name }}
