# roles/sys-env/templates/etc/update-motd.d/70-tb-grafana.j2
#!/bin/bash
echo ""
echo "Grafana url:"
echo "https://grafana.terminalboom.com/d/StDmxa3Zz/seeders-dashboard?orgId=2&refresh=10s&var-server=`hostname`"
